# contrib/docker/00-base.fedora/profile.d/prompt.sh
PROMPT_COMMAND='echo -ne "\033]0;${USER}@${HOSTNAME}: ${PWD}\007"'
