# meta/recipes/elito-base-files/elito-base-files/var-cache.automount
[Unit]
Before=local-fs.target

[Automount]
Where=/var/cache
