# meta/recipes/elito-base-files/elito-systemd-conf/device-touchscreen.target
[Unit]
Description=Touchscreen Target
Wants=dev-input-event\x2dtouchscreen0.device
After=dev-input-event\x2dtouchscreen0.device
