# contrib/docker/00-base.ubuntu/profile.d/ccache.sh
CCACHE_DIR=/opt/bld/ccache/$UID
export CCACHE_DIR
